# .config/systemd/user/clipperwatch.timer
[Unit]
Description=Run ClipperWatch daily and on boot

[Timer]
# Run 1 minute after boot
OnActiveSec=1min
# Run every 24 hours
OnUnitActiveSec=24h
# Ensure it runs if the computer was off during the scheduled time
Persistent=true

[Install]
WantedBy=timers.target
